# sadmin-mcp.service
[Unit]
Description=SADMIN 2025 MCP REST API Server
After=network.target

[Service]
Type=simple
User=buslogic
WorkingDirectory=/home/buslogic/sadmin2025-mcp-server
Environment="NODE_ENV=production"
Environment="ENVIRONMENT=production"
Environment="SADMIN_API_URL=https://sadmin2025api.ticketing.rs/api/claude"
Environment="SADMIN_API_KEY=claude_DoGQQYPIAzk2uZMwQ24o4cOMMyLsFw72"
Environment="SADMIN_PROJECT_ID=be92301b-8118-4651-a252-a996339ead2e"
Environment="MCP_SERVER_PORT=3010"
ExecStart=/usr/bin/node /home/buslogic/sadmin2025-mcp-server/dist/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target